# scripts/cloud-config/centos/setup-monitor.socket
[Socket]
ListenStream={{setupMonitorPort}}
Accept=true
NoDelay=true
